# snapkey/snapkey.service
[Unit]
Description=SnapKey daemon
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/env python3 %h/.local/share/snapkey/src/daemon.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
